# systemd/gaztime-pod.service
[Unit]
Description=Gaztime Pod POS
After=network.target gaztime-api.service
Wants=gaztime-api.service

[Service]
Type=simple
User=hein
WorkingDirectory=/home/hein/apps/gaztime/apps/pod
Environment=NODE_ENV=production
Environment=PORT=3337
Environment=VITE_API_URL=http://172.17.0.1:3333
ExecStart=/usr/bin/pnpm preview --port 3337 --host 0.0.0.0
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
